# gradle/systemd/syncany-gui.service
[Unit]
Description=Syncany GUI

[Service]
ExecStart=/usr/bin/syncany gui

[Install]
WantedBy=default.target
